# posts/incident-response.json
{
  "posts": [

    {
      "title":   "Investigando WinMemoryCleaner.exe — possível persistência via tarefa agendada",
      "date":    "2025-03-20",
      "tag":     "Análise",
      "excerpt": "Análise de incidente real envolvendo WinMemoryCleaner com possível persistência via Scheduled Task (T1053.005) em endpoint AzureAD. Processo de triagem e hunting no MDE.",
      "mitre":   ["T1053.005", "T1036", "T1204"],
      "url":     "posts/incident-response/winmemorycleaner-scheduled-task.html"
    }

  ]
}